# deploy/notisync-broker-v2.service
[Unit]
Description=NotiSync compact-CBOR broker (notisync-api-v2.extrawdw.net)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/notisync-v2
ExecStart=/usr/bin/java -Dktor.deployment.host=127.0.0.1 -Dktor.deployment.port=8082 -jar /opt/notisync-v2/server-all.jar
Environment=NOTISYNC_DB_PATH=/var/lib/notisync-v2/notisync.db
Environment=NOTISYNC_DEMO_CONFIG_PATH=/etc/notisync/demo-experience.json
Environment=NOTISYNC_SECURITY_ENABLED=true
Environment=NOTISYNC_INTEGRITY_REQUIRED=false
Environment=NOTISYNC_POW_DIFFICULTY=0
# Push-provider/App Check settings shared with the existing brokers. The per-instance file contains
# only v2-private values such as NOTISYNC_METRICS_PASSWORD and is installed mode 0600.
EnvironmentFile=/etc/notisync/common.env
EnvironmentFile=/etc/notisync/v2.env
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/notisync-v2
PrivateTmp=true

[Install]
WantedBy=multi-user.target
